# etc/services.d/one/run
#!/command/execlineb -P
while true
sleep 100
