# etc/services.d/one/run
#!/command/execlineb -P
echo "while true; do sleep 100; done" | sh -

