# etc/services.d/one/run
#!/command/execlineb -P
sh -c "echo 'script start'; while true; do sleep 7; done;"

