# etc/services.d/one/run
#!/command/execlineb -P
while true; do sleep 100; done

