# etc/services.d/asterisk/run
#!/command/execlineb -P
exec /usr/sbin/asterisk -f -p
